# bills.service
[Unit]
Description=Bill Tracker
After=network.target

[Service]
Type=simple
# Change these to match your user and install location
User=youruser
WorkingDirectory=/home/youruser/bill-tracker
ExecStart=/home/youruser/bill-tracker/venv/bin/python -m app.main
Restart=on-failure
RestartSec=10
# Configuration is read from the .env file in WorkingDirectory,
# or set environment variables here instead:
# Environment=NTFY_TOPIC=your-topic-here

[Install]
WantedBy=multi-user.target
